# startup-scripts/clearpath/clearpathd.service
[Unit]
After=network.service

[Service]
ExecStart=/usr/local/bin/install-clearpath.sh

[Install]
WantedBy=default.target
